# debian/prerm
#!/bin/bash

systemctl stop dhcp_listener
systemctl disable dhcp_listener
